# external/boot-chandelier.service
[Unit]
Description=Run Boot on Wifi Connection
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=~/Desktop/chandelier/scripts/boot.sh

[Install]
WantedBy=multi-user.target